# launchd/com.macos-gesture-launcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.macos-gesture-launcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/PATH/TO/macos-gesture-launcher</string>
    <string>--config</string>
    <string>/PATH/TO/config.yaml</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>/tmp/macos-gesture-launcher.err</string>
  <key>StandardOutPath</key>
  <string>/tmp/macos-gesture-launcher.out</string>
</dict>
</plist>
